# config/wvdial.service
[Unit]
Description=GSM Modem Dialer

[Service]
Type=simple
ExecStart=/usr/bin/wvdial
